# jenkins/gunicorn/ecom.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/run/ecom.sock
[Install]
WantedBy=sockets.target
